# units/logload_cp.service
[Unit]
Description=copying logged data to a non volatile path

[Service]
Type=oneshot
ExecStart=/bin/sh /usr/bin/log-load-cp-to-var.sh

